# personal-os/AIエージェント基盤/loops-registry/loops/watch-keeper/com.kitamura.watch-keeper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kitamura.watch-keeper</string>

  <!--
    draft。配置のみ・未登録（有効化は人間ゲート＝指揮官側が実施）。
    見張り番(watch.sh)が死んでいる間も5分毎に検知・通知する常設keeper（loop.md参照）。

    有効化（cpではなくsymlink。常設loopのためリポ編集がそのままlaunchd登録へ反映されるようにする）:
      ln -s '/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/watch-keeper/com.kitamura.watch-keeper.plist' \
        ~/Library/LaunchAgents/com.kitamura.watch-keeper.plist
      mkdir -p '/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/watch-keeper/output/logs'
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.kitamura.watch-keeper.plist
      launchctl enable gui/$(id -u)/com.kitamura.watch-keeper
    停止:
      launchctl bootout gui/$(id -u)/com.kitamura.watch-keeper
      rm ~/Library/LaunchAgents/com.kitamura.watch-keeper.plist
    手動実行（通知/alerts追記が実際に起きる）:
      cd <このフォルダ> && scripts/keeper.sh
    RunAtLoad は付けない（daily-digest/inbox-patrolと同じ方針）。keeper.shは読み取り専用判定＋
    通知/追記のみで非破壊・冪等（loop.md §冪等性）。
  -->

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd '/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/watch-keeper' &amp;&amp; scripts/keeper.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/watch-keeper</string>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/kitamuranaohiro</string>
    <key>PATH</key>
    <string>/Users/kitamuranaohiro/.local/bin:/Users/kitamuranaohiro/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>Asia/Tokyo</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/watch-keeper/output/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/watch-keeper/output/logs/launchd.err.log</string>
</dict>
</plist>
