# personal-os/AIエージェント基盤/loops-registry/loops/daily-notion-sync/com.kitamura.daily-notion-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kitamura.daily-notion-sync</string>

  <!--
    30秒ごと（StartInterval 30）に当日デイリーの「動いているエージェント」「終わったこと」の
    変化を検知し、変化がある時だけNotionの表A/表B（session-table.sh）を更新するsync loop
    （2026-07-06デイリーNotion表反映plan.md）。無変化ならAPI呼び出しゼロで即exit 0。

    draft。置くだけ・未ロード（有効化は人間ゲート＝GLOBAL_AGENTS.md §7の破壊操作ゲート対象）。
    有効化する場合（symlink方式。既存draft plistと同じ実運用パターン）:
      ln -s '/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/daily-notion-sync/com.kitamura.daily-notion-sync.plist' \
        ~/Library/LaunchAgents/com.kitamura.daily-notion-sync.plist
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.kitamura.daily-notion-sync.plist
      launchctl enable gui/$(id -u)/com.kitamura.daily-notion-sync
    停止:
      launchctl bootout gui/$(id -u)/com.kitamura.daily-notion-sync
    手動実行: cd <このフォルダ> && scripts/sync.sh

    RunAtLoad は付けない（起動直後の連打を避ける。既存draft plistと同じ方針）。
    sync.sh 自体が多重起動防止ロック（state/notion-session-table-sync.lock・mkdirベース・
    300秒でstale自己修復）を持つため、前回runが30秒を超えて残っていても二重実行しない。
  -->

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd '/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/daily-notion-sync' &amp;&amp; scripts/sync.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/daily-notion-sync</string>

  <key>StartInterval</key>
  <integer>30</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/kitamuranaohiro</string>
    <key>PATH</key>
    <string>/Users/kitamuranaohiro/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>Asia/Tokyo</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/daily-notion-sync/output/logs/sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/daily-notion-sync/output/logs/sync.err.log</string>
</dict>
</plist>
