# personal-os/AIエージェント基盤/loops-registry/loops/ai-jobs-dispatcher/com.kitamura.ai-jobs-dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kitamura.ai-jobs-dispatcher</string>

  <!--
    置くだけ・未ロード（有効化は人間ゲート）。有効化する場合:
      cp com.kitamura.ai-jobs-dispatcher.plist ~/Library/LaunchAgents/
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.kitamura.ai-jobs-dispatcher.plist
      launchctl enable gui/$(id -u)/com.kitamura.ai-jobs-dispatcher
    停止:
      launchctl bootout gui/$(id -u)/com.kitamura.ai-jobs-dispatcher
  -->

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd '/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/ai-jobs-dispatcher' &amp;&amp; '/Users/kitamuranaohiro/.npm-global/bin/tsx' scripts/dispatcher.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/ai-jobs-dispatcher</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/kitamuranaohiro</string>
    <key>PATH</key>
    <string>/Users/kitamuranaohiro/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>Asia/Tokyo</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/ai-jobs-dispatcher/output/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/ai-jobs-dispatcher/output/logs/launchd.err.log</string>
</dict>
</plist>
