# personal-os/AIエージェント基盤/loops-registry/loops/renderer/com.kitamura.lanes-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kitamura.lanes-sync</string>

  <!--
    毎分（StartInterval 60）cockpitレーン（orca worktree ps）の変化を検知し、変化がある時だけ
    ローカル auto:board-now とNotion「レーン実況」DB（N3b・notion-lanes.sh）を更新するsync loop
    （統合program plan.md 方針5c）。無変化ならAPI呼び出しゼロで即exit 0。

    draft。置くだけ・未ロード（有効化は人間ゲート＝運用契約§2・§7の破壊操作ゲート対象）。
    有効化する場合（symlink方式。inbox-patrolと同じ実運用パターン）:
      ln -s '/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/renderer/com.kitamura.lanes-sync.plist' \
        ~/Library/LaunchAgents/com.kitamura.lanes-sync.plist
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.kitamura.lanes-sync.plist
      launchctl enable gui/$(id -u)/com.kitamura.lanes-sync
    停止:
      launchctl bootout gui/$(id -u)/com.kitamura.lanes-sync
    手動実行: cd <このフォルダ> && scripts/lanes-sync.sh

    RunAtLoad は付けない（起動直後の連打を避ける。既存draft plistと同じ方針）。
    lanes-sync.sh 自体が多重起動防止ロック（state/notion-lanes-sync.lock・mkdirベース・
    300秒でstale自己修復）を持つため、前回runが60秒を超えて残っていても二重実行しない。
  -->

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd '/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/renderer' &amp;&amp; scripts/lanes-sync.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/renderer</string>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/kitamuranaohiro</string>
    <key>PATH</key>
    <string>/Users/kitamuranaohiro/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>Asia/Tokyo</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/renderer/output/logs/lanes-sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kitamuranaohiro/Private/personal-os/AIエージェント基盤/loops-registry/loops/renderer/output/logs/lanes-sync.err.log</string>
</dict>
</plist>
